# deployments/linux/systemd/resqauto-hello.service
[Unit]
Description=ResQAuto Hello world (one-shot, Docker)
Wants=docker.service network-online.target
After=docker.service network-online.target

[Service]
Type=oneshot
# Optional: file with environment like TZ, extra vars
EnvironmentFile=/etc/default/resqauto-hello
# Run your tagged image; adjust the image/tag if you publish it
ExecStart=/usr/bin/docker run --rm --env TZ=${TZ:-Europe/Stockholm} resqauto:hello
# Avoid lingering containers if something goes wrong
TimeoutStartSec=300